# systemd/meeting-recorder-summarizer.container
[Unit]
Description=Meeting Recorder Summarizer Service
Documentation=https://github.com/tu-usuario/meeting-recorder-summarizer
After=pipewire.service

[Container]
Image=localhost/meeting-recorder-summarizer:latest
ContainerName=meeting-recorder-summarizer

# Volumen de datos: grabaciones, transcripciones, resúmenes y modelos de Whisper
Volume=%h/meeting-recorder-summarizer-data:/data

# Montar el runtime dir completo del host (contiene sockets de PipeWire y PulseAudio)
# %t = $XDG_RUNTIME_DIR del host (ej: /run/user/1000)
Volume=%t:/run/user/host

# Variables de entorno del servicio
# XDG_RUNTIME_DIR apunta al runtime dir montado del host
Environment=XDG_RUNTIME_DIR=/run/user/host
Environment=PULSE_SERVER=unix:/run/user/host/pulse/native
Environment=WHISPER_MODEL=small
Environment=ENABLE_CUDA=false
Environment=STORAGE_PATH=/data

# OPENAI_API_KEY se lee desde un archivo de entorno del usuario
EnvironmentFile=%h/.config/meeting-recorder-summarizer.env

# Comando a ejecutar (serve = modo daemon esperando comandos)
Exec=serve

# Mantener el UID del usuario del host dentro del contenedor
UserNS=keep-id

# Deshabilitar SELinux labeling para permitir acceso a los sockets
SecurityLabelDisable=true

[Service]
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60

[Install]
WantedBy=default.target
